# crates/service/macos/nz.ophelia.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>nz.ophelia.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>__OPHELIA_SERVICE_BINARY__</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>nz.ophelia.service</key>
        <true/>
    </dict>
    <key>KeepAlive</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>__OPHELIA_LOG_DIR__/ophelia-service.out.log</string>
    <key>StandardErrorPath</key>
    <string>__OPHELIA_LOG_DIR__/ophelia-service.err.log</string>
</dict>
</plist>
